# root/etc/services.d/supysonic-watcher/run
#!/usr/bin/with-contenv sh

umask 022

exec s6-setuidgid appuser \
   supysonic-watcher

